# docker-script3/supervisord.conf
[program:sonoff-script3]
priority=10
directory=/scripts/problema4/
command=python3 script3.py
user=root
autostart=true
autorestart=true