# scripts/ejemplos/docker-script3/supervisord.conf
[program:sonoffscript3]
priority=10
directory=/scripts/problema4/
command=python3 script3.py
autostart=true
autorestart=true
user=root