# docker-script2/supervisord.conf
[program:sonoffscript2]
priority=10
directory=/scripts/problema4/
command=python3 script2.py
autostart=true
autorestart=true
user=root