# scripts/proyecto/3-docker/supervisord.conf
[program:cerradura]
priority=10
directory=../programa/
command=python3 publisher-docker.py JANO
autostart=true
autorestart=true
user=root