# scripts/Cerradura/docker/supervisord.conf
[program:cerradura]
priority=10
directory=../programa/
command=python3 publisher.py JANO
autostart=true
autorestart=true
user=root